# pulse/99-hdmi-pulse.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/switch-audio.sh"
